# overlay/etc/cron.hourly/wireguard-profiles-delexpired
#!/bin/bash -e
# delete downloaded wireguard client profiles

su www-data -s /bin/bash -c /var/www/wireguard/bin/delexpired
